# sysutils/bsdstats/files/bsdstats.in
#!/sbin/openrc-run
#

name=bsdstats

command=%%PREFIX%%/etc/periodic/monthly/300.statistics
command_args="-nodelay >/dev/null 2>/dev/null &"
